# tools/config/cabot-config.service
[Unit]
Description = Change Configuration For CaBot

[Service]
ExecStart = /opt/cabot/tools/cabot_config.sh
Restart = no
WorkingDirectory = /opt/cabot
Type = oneshot

[Install]
WantedBy = default.target